# scripts/com.01v96remote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.01v96remote.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/YourUsername/01v96-remote/server.js</string>
        <string>midi</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/YourUsername/01v96-remote</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/YourUsername/01v96-remote/logs/server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/YourUsername/01v96-remote/logs/server-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
